# www-apps/ttyd/files/ttyd.initd
#!/sbin/openrc-run
# shellcheck shell=sh
# Copyright 2023 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

name="ttyd daemon"
command="/usr/bin/ttyd"
command_args="${TTYD_OPTS}"
command_background="true"
pidfile="/run/ttyd.pid"

depend() {
	use logger net
}
